# coffee.service
[Unit]
Description=Coffee Tracker API
After=network.target

[Service]
Type=simple
User=coffee
Group=coffee
WorkingDirectory=/opt/coffee
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
EnvironmentFile=/opt/coffee/.env
ExecStart=/usr/bin/env node /opt/coffee/dist/index.js
Restart=always
RestartSec=5

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coffee

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/coffee/data /opt/coffee/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
